# gentoo/etc/profile.d/arm-andorid.sh
#!/bin/bash

export PATH="$PATH:/opt/arm-android15/bin:/opt/arm-linux/bin"
