# gentoo/etc/profile.d/colorgcc.sh
#!/bin/bash
###############################################################################
# Author : liyunteng
# Email : li_yunteng@163.com
# Created Time : 2015-06-12 19:35
# Filename : colorgcc.sh
# Description : 
###############################################################################
#export PATH="/usr/lib/colorgcc/bin:${PATH}"
