# gentoo/etc/profile.d/go.sh
#!/bin/sh
export GOPATH=/home/lyt/work/go
export PATH=$PATH:${GOPATH}/bin
